# etc/supervisord.d/shortener.ini
[program:aiohttp]
numprocs = 6
numprocs_start = 1
process_name = shortener_%(process_num)s
environment=PYTHONPATH=/var/www/shortener
command=python3.6 -m shortener.main --port=808%(process_num)d
user=nginx
autostart=true
autorestart=true
