# etc/supervisord.d/shortener.ini
[program:shortener]
numprocs = 3
numprocs_start = 1
process_name = shortener_%(process_num)s
command=python3 -m shortener -c /etc/shortener.yml -p 808%(process_num)d
stdout_logfile=/var/log/shortener.log
redirect_stderr=true
autostart=true
autorestart=true
